# dbus-org.sailfishos.usermanager.service
[Unit]
Description=Sailfish user manager service

[Service]
Type=dbus
BusName=org.sailfishos.usermanager
ExecStart=/usr/bin/user-managerd
